# systemd/user/openclaw-gateway.service.d/10-moltbook-env.conf
[Service]
EnvironmentFile={{ OPENCLAW_HOME }}/moltbook.env
