# systemd/user/moltbook-relay.service
[Unit]
Description=Relay moltbook staging files to sanitizer input

[Service]
Type=oneshot
ExecStart={{ OPENCLAW_WORKSPACE }}/scripts/moltbook-relay.sh
